# ubuntu18.04/overlay-adv/etc/udev/rules.d/77-mm-quectel-port-types.rules
SUBSYSTEMS=="usb", ACTION=="add" ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0800", RUN+="/usr/bin/adv-mmscan-hook.sh"

